# dash-deploy/rundashboard.service
[Unit]
Description=Run Dash Dashboard.

[Service]
Type=simple
ExecStart=/bin/bash /usr/bin/rundashboard.sh

[Install]
WantedBy=multi-user.target
